# ETC/profile.d/aws.sh
export AWS_DEFAULT_PROFILE=default
complete -C /usr/local/bin/aws_completer aws
