# ETC/profile.d/env.sh
export JAVA_HOME=`/usr/libexec/java_home -v 1.8`
export GRADLE_HOME=/usr/local/gradle
export HOMEBREW_GITHUB_API_TOKEN=6063b49b02621de5d677a4c985ff9e0a648ade6e

bind 'set match-hidden-files off' # Stop autocomplete considering dot files

