# ETC/profile.d/env.sh
export JAVA_HOME=`/usr/libexec/java_home -v 1.8`
export GRADLE_HOME=/usr/local/gradle
export HOMEBREW_GITHUB_API_TOKEN=7f52a33d921712c32b288d7c9b62a807e39feb36
bind 'set match-hidden-files off' # Stop autocomplete considering dot files
# Use system_profiler SPDisplaysDataType |grep Resolution to get resolution.  I chose one
# Use X_WIDTH and X_HEIGHT with display_full, an alias with args for Imagemagicks display
# Height is just a bit less to allow for menu bar
export X_WIDITH=1920
export X_HEIGHT=1040
