# ETC/profile.d/paths.sh
export PATH=/usr/local/bin:$JAVA_HOME/bin:/usr/local/sbin:$PATH
export PATH=/Applications/VirtualBox.app/Contents/Resources/VirtualBoxVM.app/Contents/MacOS:$PATH

